# btc-logger.service
[Unit]
Description=BTC Price Logger - Always running logger service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/polymarket-bot
ExecStart=/usr/bin/python3 /home/ubuntu/polymarket-bot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Auto-restart settings
StartLimitInterval=0
StartLimitBurst=0

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
